# Roboshop-ANSIBLE/ansible-roles/roles/shipping/files/shipping.service
[Unit]
Description=Shipping Service

[Service]
User=roboshop
Environment=CART_ENDPOINT=cart.joiningindevops.online:8080
Environment=DB_HOST=mysql.joiningindevops.online
ExecStart=/bin/java -jar /app/shipping.jar
SyslogIdentifier=shipping

[Install]
WantedBy=multi-user.target